# crontab
03 18 * * * phantomjs $HOME/phantom/attendance.js
